# services/rovi-nav.service
[Unit]
Description=ROVI stack: nav
Requires=rovi-gateway.service
After=rovi-gateway.service
PartOf=rovi-gateway.service

[Service]
Type=simple
User=rovi
WorkingDirectory=%h/dev/Roblibs/rovi_ros_ws
EnvironmentFile=%h/dev/Roblibs/rovi_ros_ws/.env.robot
RuntimeDirectory=rovi
ExecStart=%h/dev/Roblibs/rovi_ros_ws/services/run_stack.sh nav
KillSignal=SIGINT
TimeoutStopSec=20
Restart=no

[Install]
WantedBy=multi-user.target

